# archiso-systemd/airootfs/etc/rc.local
#!/bin/sh

cd /usr/share
rm -rf kb iana-etc zoneinfo file info man gettext doc gtk-doc hwdata bash-completion licenses gir-1.0 ca-certificates zoneinfo-leaps audit readline vala icu gettext nano fish factory aclocal 
find locale ! -name C | xargs rm -rf
find i18n/locales ! -name C | xargs rm -rf
cd /usr/bin
rm -rf fsck.ext2 fsck.ext3 locale-gen localectl localedef locale
